# edge/systemd/alsa-restore.service.j2
[Unit]
Description=Restore ALSA mixer settings
After=sound.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/alsactl -f /var/lib/alsa/asound.state restore 1

[Install]
WantedBy=multi-user.target 