# edge/systemd/darkice.service.j2
[Unit]
Description=DarkIce Live Audio Streamer
Documentation=man:darkice(1)
After=network-online.target sound.target
Wants=network-online.target
Requires=sound.target

[Service]
Type=simple
ExecStart=/usr/bin/darkice -c {{ sculpture_dir }}/darkice.cfg
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
User=pi
Group=audio

# Environment variables
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment=SCULPTURE_ID={{ id }}
Environment=CONTROL_HOST={{ control_host }}

# Security settings
NoNewPrivileges=true
PrivateTmp=true
AmbientCapabilities=CAP_SYS_NICE

[Install]
WantedBy=multi-user.target 