# edge/systemd/player-live.service.j2
[Unit]
Description=Sculpture Live Audio Player
After=network-online.target sound.target
Wants=network-online.target
Requires=sound.target

[Service]
Type=simple
ExecStart=/usr/bin/mpv --no-video --audio-device=pulse/sculpture_sink --cache=yes --cache-secs=10 --demuxer-max-bytes=5M --audio-buffer=2 http://{{ control_host }}:8000/mix-for-{{ id }}.ogg
Restart=always
RestartSec=5
User=pi
Group=audio

# Environment variables
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment=SCULPTURE_ID={{ id }}
Environment=CONTROL_HOST={{ control_host }}

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target 