# edge/systemd/player-live.service.j2
[Unit]
Description=Sculpture Live Audio Player
After=network-online.target sound.target
Wants=network-online.target
Requires=sound.target

[Service]
Type=simple
ExecStart=/usr/bin/mpv --no-video --audio-device={{ mpv_audio_device }} --audio-samplerate={{ audio_sample_rate }} --audio-format={{ mpv_audio_format }} --cache=yes --cache-secs=60 --demuxer-max-bytes=20M --audio-buffer=5 --msg-level=all=v --log-file=/tmp/mpv-live-{{ id }}.log --input-ipc-server=/tmp/mpvsocket-{{ id }} http://{{ control_host }}:8000/mix-for-{{ id }}.ogg
Restart=always
RestartSec=5
User=pi
Group=audio
Nice=-10
SuccessExitStatus=4

# Environment variables
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment=SCULPTURE_ID={{ id }}
Environment=CONTROL_HOST={{ control_host }}

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target 