# server/liquidsoap/mqtt_to_telnet_bridge.service
[Unit]
Description=MQTT to Liquidsoap Telnet Bridge
After=network.target

[Service]
Type=simple
User=unix_user
WorkingDirectory=/opt/sculpture-system/liquidsoap
ExecStart=/usr/bin/python3 /opt/sculpture-system/liquidsoap/mqtt_to_telnet_bridge.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target 
