# docker/pi-simulator/systemd-services/pi-agent.service
[Unit]
Description=Sculpture Pi MQTT Agent (Simulator)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/sculpture-system/pi-agent.py
Restart=always
RestartSec=5
User=pi
Group=pi

# Environment variables
Environment=SCULPTURE_ID=1
Environment=CONTROL_HOST=host.docker.internal

# Security settings
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target 